# conf/puppet/cron.d/03-update
#!/bin/bash
sleep ${1:-$(($RANDOM%10))}
apt-get update -o Acquire::http::Timeout="15";
